# etc/local.d/macchanger.start
#!/usr/bin/env bash
interface=wlp2s0
ip link set $interface down
/usr/bin/macchanger -r $interface
ip link set $interface up
